# dot_config/systemd/user/user-session.target
[Unit]
Description=User sesion target
AllowIsolate=true
