# dot_config/systemd/user/user-timers.target
[Unit]
Description=User level timers target
AllowIsolate=true

[Install]
WantedBy=user-session.target
